# lib/systemd/system/wlan-assistant.service
[Unit]
Description=WLAN Assistant
After=droid-hal-init.service
Before=wlan-module-load.service

[Service]
Type=basic
Restart=always
ExecStart=/vendor/bin/wlan_assistant

[Install]
WantedBy=multi-user.target
